# verticals/bitnin/services/bitnin_observability/systemd/bitnin-shadow.timer
[Unit]
Description=Run BitNin Shadow Supervisor periodically

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
